# crypt-admin/init-phoenix/config/etc--systemd--system--daily.service
[Service]
ExecStart=bash /root/crypt/program-files/daily-service/daily-service.sh