# crypt-admin/init-phoenix/config/etc--systemd--system--minutely.service
[Service]
ExecStart=bash /root/crypt/program-files/minutely-service/minutely-service.sh