# crypt-admin/init-phoenix/config/etc--systemd--system--minutely.timer.template
[Timer]
OnActiveSec=300
OnUnitInactiveSec=300

[Install]
WantedBy=multi-user.target